# admin/fixvox-web/systemd/fixvox-release-broker.service
[Unit]
Description=Fixvox gated Git and release broker
After=network-online.target fixvox-admin-deploy-helper.service
Wants=network-online.target

[Service]
Type=simple
User=fixvox-release
Group=fixvox-release-broker
RuntimeDirectory=fixvox-release
RuntimeDirectoryMode=0770
RuntimeDirectoryPreserve=yes
Environment=PI_CHAT_RELEASE_BROKER_ENABLED=1
Environment=PI_CHAT_RELEASE_BROKER_SOCKET=/run/fixvox-release/release.sock
Environment=PI_CHAT_RELEASE_CONFIG=/etc/fixvox-release/release.json
Environment=PI_CHAT_RELEASE_JOURNAL=/var/lib/fixvox-release/audit/release.jsonl
Environment=PI_ADMIN_DEPLOY_SOCKET=/run/fixvox-release/admin-deploy.sock
ExecStart=/usr/bin/node /opt/fixvox-agent/runtime/pi-release-service.mjs
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/fixvox-workspace/repos/dictation-tauri /var/lib/fixvox-release /run/fixvox-release
ReadOnlyPaths=/etc/fixvox-release /opt/fixvox-agent/runtime
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
