# admin/fixvox-web/systemd/fixvox-constelaciones-read-broker.service
[Unit]
Description=Fixvox redacted Constelaciones read broker
After=network.target

[Service]
Type=simple
User=root
Group=fixvox-agent-broker
RuntimeDirectory=fixvox-agent
RuntimeDirectoryMode=0770
RuntimeDirectoryPreserve=yes
Environment=PI_CHAT_CONSTELACIONES_SOCKET=/run/fixvox-agent/constelaciones-read.sock
Environment=CONSTELACIONES_READ_DB=/var/lib/docker/volumes/i103nz67zizmtz5m2rt5k9e3-constelaciones-data/_data/constelaciones.sqlite
ExecStart=/usr/bin/node /opt/fixvox-agent/runtime/constelaciones-read-broker.mjs
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/var/lib/docker/volumes/i103nz67zizmtz5m2rt5k9e3-constelaciones-data/_data/constelaciones.sqlite
ReadWritePaths=/run/fixvox-agent
RestrictAddressFamilies=AF_UNIX
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
